# etc/cron.d/parentrol
# /etc/cron.d/anacron: crontab entries for the anacron package

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
TOOLSDIR=/usr/share/parentrol
ETCDIR=/etc
LOCKDIR=/var/lock

* * * * * root $TOOLSDIR/check.sh
@reboot root $TOOLSDIR/check.sh
